# dotfiles/etc/udev/rules.d/10-local.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/systemctl start monitors@james.service"
